# com.chesshacks.overnight-training.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.chesshacks.overnight-training</string>
	
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/0591a/development/projects/ChessHacks/ChessHacks-CLI/run_overnight_training.sh</string>
		<string>--depth</string>
		<string>15</string>
		<string>--hours</string>
		<string>8</string>
	</array>
	
	<!-- Run at 10 PM every day -->
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>22</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	
	<!-- Allow the job to run for up to 10 hours -->
	<key>TimeOut</key>
	<integer>36000</integer>
	
	<!-- Keep job running if it crashes -->
	<key>KeepAlive</key>
	<false/>
	
	<!-- Redirect output -->
	<key>StandardOutPath</key>
	<string>/Users/0591a/development/projects/ChessHacks/ChessHacks-CLI/logs/launchd.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/0591a/development/projects/ChessHacks/ChessHacks-CLI/logs/launchd_error.log</string>
	
	<!-- Working directory -->
	<key>WorkingDirectory</key>
	<string>/Users/0591a/development/projects/ChessHacks/ChessHacks-CLI</string>
	
	<!-- User environment -->
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
</dict>
</plist>
